# setup-stuffs/server-stuff/wastebin.service
[Unit]
Description=Wastebin

[Service]
ExecStart=/home/homeless/wastebin/wastebin.sh

[Install]
WantedBy=multi-user.target
